# KernovaGuestAgent/app.kernova.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>app.kernova.agent</string>
	<!-- Tells Login Items & Extensions ("Allow in the Background") which app
	     bundle owns this job, so it shows Kernova Guest Agent.app's icon + display
	     name instead of a generic icon. Required because we install this plist
	     into ~/Library/LaunchAgents directly rather than via SMAppService. -->
	<key>AssociatedBundleIdentifiers</key>
	<array>
		<string>app.kernova.agent</string>
	</array>
	<key>ProgramArguments</key>
	<array>
		<string>__INSTALL_DIR__/Kernova Guest Agent.app/Contents/MacOS/KernovaGuestAgent</string>
	</array>
	<!-- Vends the agent↔File Provider extension XPC relay (#376). launchd owns the
	     Mach service and hands it to the agent's NSXPCListener; the name is app-group
	     prefixed so the sandboxed extension is allowed to look it up. -->
	<key>MachServices</key>
	<dict>
		<key>8MT4P4GZL2.app.kernova.relay</key>
		<true/>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>
	<key>ThrottleInterval</key>
	<integer>10</integer>
</dict>
</plist>
